# infra/roles/services/reverse-proxy/firewall-tunnel/files/proxy-http-to-local-port.socket
[Unit]
Wants=sockets.target
PartOf=proxy-http-to-local-port.service

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
