# infra/roles/services/reverse-proxy/firewall-tunnel/files/proxy-https-to-local-port.socket
[Unit]
Wants=sockets.target
PartOf=proxy-https-to-local-port.service

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target
